# com.nickflorez.png-organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nickflorez.png-organizer</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nickflorez/ai-assistant-onboarding/claude-png-organizer.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/nickflorez/ai-assistant-onboarding/png-organizer.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/nickflorez/ai-assistant-onboarding/png-organizer-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/nickflorez/ai-assistant-onboarding</string>
</dict>
</plist>